# supervisor/courier.conf
[program:courier]
directory=/opt/courier/src
command=python courier.py 80
